# inittab
#/dev/console::respawn:/bin/login
/dev/console::sysinit:/bin/linitd --unix /tmp/linit
